# debian/housing-data.upstart
description "Housing Data service"
author "mygov.scot team"

start on runlevel [2345]
stop on runlevel [!2345]

setuid housing-data
chdir /opt/housing-data

script
    . /etc/profile

    if [ -x /usr/local/bin/am-i-live ] && ! /usr/local/bin/am-i-live; then
      exit
    el
      exec node src/run.js >> /var/log/housing-data/housing-data.log 2>&1
    fi
end script
