# debian/housing-data.service
[Unit]
Description=Housing Data Service
After=network.target

[Service]
User=housing-data
WorkingDirectory=/opt/housing-data
ExecStart=/opt/housing-data/bin/housing-data
Restart=on-failure

[Install]
WantedBy=multi-user.target
