# debian/housing-data.service
[Unit]
Description=Housing Data Service
After=network.target

[Service]
User=housing-data
WorkingDirectory=/opt/housing-data
ExecStart=/bin/bash -c "source /etc/profile; exec /usr/bin/node src/server.js >> /var/log/housing-data/housing-data.log 2>&1"
Restart=on-failure
[Install]
WantedBy=multi-user.target
